# nuxt copy/rootfs/etc/s6-overlay/s6-rc.d/app/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: AdGuard Home
# Runs the AdGuard Home server
# ==============================================================================
declare -a options

bashio::log.info "Starting Nuxt App server..."

# RUN AdGuard Home server
cd /opt/app || exit
if ! bashio::fs.file_exists ".output/server/index.mjs"; then
    exit
fi
exec node .output/server/index.mjs
